# config/supervisor-telcoflo.example.conf
; Copy to server: sudo cp config/supervisor-telcoflo.example.conf /etc/supervisor/conf.d/queue-worker.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update
;
; 3-pool layout for 32GB RAM. Scale by increasing numprocs (e.g. 12–15 per pool).
; Each worker uses ~50–128MB. --max-jobs=1000 recycles workers to limit memory creep (--max-memory requires Laravel 8+).
; Ensure "default" is in one pool so no jobs are left unprocessed.

# POOL 1: Billing & high priority
[program:worker-billing]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/telcoflo_studio/artisan queue:work database --queue=billing,high --sleep=1 --tries=3 --max-jobs=1000
directory=/var/www/telcoflo_studio
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=12
redirect_stderr=true
stdout_logfile=/var/log/telcoflo/worker-billing.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopwaitsecs=3600

# POOL 2: SMS campaigns
[program:worker-sms]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/telcoflo_studio/artisan queue:work database --queue=sms --sleep=1 --tries=3 --max-jobs=1000
directory=/var/www/telcoflo_studio
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=12
redirect_stderr=true
stdout_logfile=/var/log/telcoflo/worker-sms.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopwaitsecs=3600

# POOL 3: Low priority & default (reports, delivery status, and any uncategorized jobs)
[program:worker-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/telcoflo_studio/artisan queue:work database --queue=low,default --sleep=1 --tries=3 --max-jobs=1000
directory=/var/www/telcoflo_studio
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/log/telcoflo/worker-low.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stopwaitsecs=3600
